# trance-daemon/debian/prerm
#!/bin/sh
# prerm for trance — only stop/disable on full remove.
# On upgrade, leave the daemon running; postinst restarts onto the new binary.
# (Stopping on upgrade races with packaging order and leaves the service dead.)
set -u

for_each_user_session() {
    _cb="$1"
    command -v loginctl >/dev/null 2>&1 || return 0
    command -v systemctl >/dev/null 2>&1 || return 0
    loginctl list-users --no-legend 2>/dev/null | while read -r uid user _rest; do
        case "$uid" in ''|*[!0-9]*) continue ;; esac
        [ -n "$user" ] || continue
        [ -d "/run/user/$uid" ] || continue
        [ -S "/run/user/$uid/bus" ] || continue
        "$_cb" "$uid" "$user" || true
    done
}

_user_systemctl() {
    _uid="$1"; _user="$2"; shift 2
    if command -v runuser >/dev/null 2>&1; then
        runuser -u "$_user" -- env \
            XDG_RUNTIME_DIR="/run/user/$_uid" \
            DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$_uid/bus" \
            systemctl --user "$@" 2>/dev/null && return 0
    fi
    systemctl --user --machine="${_user}@" "$@" 2>/dev/null || true
}

try_stop_trance() {
    echo "-> stop idle-daemon for $2 (remove)"
    _user_systemctl "$1" "$2" stop idle-daemon.service || true
}

case "${1:-}" in
    remove|deconfigure)
        for_each_user_session try_stop_trance
        ;;
    upgrade|failed-upgrade)
        # Leave service running across upgrades.
        ;;
esac

if [ "${1:-}" = "remove" ]; then
    if command -v deb-systemd-helper >/dev/null 2>&1; then
        deb-systemd-helper --user disable idle-daemon.service >/dev/null 2>&1 || true
    fi
fi

exit 0
