# trance-daemon/debian/prerm
#!/bin/sh
# prerm script for trance-daemon
set -e

if [ "$1" = "remove" ] || [ "$1" = "upgrade" ]; then
    echo "=========================================================="
    echo "  Stopping Trance Screensaver Daemon Services..."
    echo "=========================================================="

    # Stop the service for all currently logged-in users using systemd's user manager API
    if command -v loginctl >/dev/null && command -v systemctl >/dev/null; then
        loginctl list-users --no-legend | while read -r uid user linger state; do
            if [ -n "$user" ]; then
                echo "-> Stopping service for: $user (UID: $uid)"
                # Try runuser first (does not require systemd-machined)
                if command -v runuser >/dev/null; then
                    if runuser -u "$user" -- env XDG_RUNTIME_DIR="/run/user/$uid" DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$uid/bus" systemctl --user stop trance-daemon.service; then
                        echo "-> Service stopped successfully for $user via runuser"
                        continue
                    fi
                fi
                # Fallback to --machine machine-control if runuser is unavailable or fails
                systemctl --user --machine="$user@" stop trance-daemon.service || true
            fi
        done
    fi
    echo "  Services stopped successfully."
    echo "=========================================================="
fi

if [ "$1" = "remove" ]; then
    # Disable the user service globally using standard deb-systemd-helper
    if command -v deb-systemd-helper >/dev/null; then
        deb-systemd-helper --user disable trance-daemon.service || true
    fi
fi

exit 0
